# systemd/user/metrics-idle-watcher.service
[Unit]
Description=Idle-triggered Chromium metrics kiosk watcher
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=%h/.local/bin/metrics-idle-watcher.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=metrics-idle-watcher
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
